# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
user=root

[program:openclaw]
command=openclaw gateway
directory=/home/node/.openclaw
user=node
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=10
stopsignal=TERM
stdout_logfile=/var/log/supervisor/openclaw.log
stderr_logfile=/var/log/supervisor/openclaw-error.log
environment=HOME="/home/node",USER="node"
